# require/libs/super_snort.conf
[inet_http_server]
port = 127.0.0.1:9001

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://127.0.0.1:9001

[program:snort_parser]
command=python3 /root/libs/snort_parser.py

[program:snort_api]
command=python3 /root/libs/snort_api.py

[program:snort]
command=sh -c "/usr/local/bin/snort -i $NETINT -A unsock -c /etc/snort/snort.conf"
autorestart=true
logfile=/dev/stdout
logfile_maxbytes=0
